# target/linux/qualcommax/ipq60xx/base-files/etc/init.d/fix_usb
#!/bin/sh /etc/rc.common

# CONFIG_BUSYBOX_CUSTOM=y CONFIG_BUSYBOX_CONFIG_DEVMEM=y CONFIG_KERNEL_DEVMEM=y

START=99

boot() {
	case $(board_name) in
	jdcloud,re-cs-07|\
	jdcloud,ax6600)
		devmem 0x1016000 32 0x201
		devmem 0x1016004 32 0x2
		logger "AX6600 USB Open"
	;;
	esac
}
