# target/linux/qualcommax/ipq60xx/base-files/etc/init.d/usb_on
#!/bin/sh /etc/rc.common
# Used to deal with ax6600 USB3.0 can not automatically start power supply problem
START=99

boot() {
    case $(board_name) in
    jdcloud,ax6600)
        MEM_ADDRESS_1="0x1016000"
        MEM_VALUE_1="0x201"
        MEM_ADDRESS_2="0x1016004"
        MEM_VALUE_2="0x2"

        devmem $MEM_ADDRESS_1 32 $MEM_VALUE_1
        devmem $MEM_ADDRESS_2 32 $MEM_VALUE_2
    ;;
*)
    ;;
    esac
}

exit 0